# debian/lavapdu-daemon.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/lavapdu-listen" ]; then
	invoke-rc.d lavapdu-listen stop || exit $?
fi

if [ -x "/etc/init.d/lavapdu-runner" ]; then
	invoke-rc.d lavapdu-runner stop || exit $?
fi

#DEBHELPER#

